# .config/systemd/user/backup@schedule.timer
[Timer]
OnCalendar=*-*-* *:00/03:00
AccuracySec=1m

