# .config/systemd/user/demo@schedule.timer
[Timer]
OnCalendar=*-*-* *:00/01:00
AccuracySec=1m
#Unit=demo@sample.service
